# scripts/systemd/npd-smee-webhook.service
[Unit]
Description=NPD GitHub Webhook Forwarder (smee.io)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=pbrown
WorkingDirectory=/home/pbrown/Novus-db
# IMPORTANT: Replace YOUR_CHANNEL with actual smee.io channel URL
# Create a channel at https://smee.io/new
ExecStart=/home/pbrown/.npm-global/bin/smee -u https://smee.io/YOUR_CHANNEL --target http://localhost:6701/api/v1/webhooks/github
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
